# os/systemd/counteros-api.service
[Unit]
Description=CounterOS local API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=counteros
WorkingDirectory=/opt/counteros/backend
EnvironmentFile=/etc/counteros/counteros.env
ExecStart=/opt/counteros/.venv/bin/python -m uvicorn app.main:app --host ${COUNTEROS_HOST} --port ${COUNTEROS_BACKEND_PORT}
Restart=on-failure
RestartSec=5

# Logs: journalctl -u counteros-api.service

[Install]
WantedBy=multi-user.target
